# infra/oracle/systemd/dpp-db-backup.service
[Unit]
Description=DPP PostgreSQL backup to OCI Object Storage
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dpp-db-backup backup
WorkingDirectory=/opt/dpp
User=root
Group=root
TimeoutStartSec=2h
