# infra/oracle/systemd/dpp-db-backup-verify.service
[Unit]
Description=Verify latest DPP PostgreSQL backup restore metadata
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dpp-db-backup verify
WorkingDirectory=/opt/dpp
User=root
Group=root
